# src/packages/fff/fff-gateway/files/etc/uci-defaults/10-gateway-network
#!/bin/sh

uci batch <<EOF
    set network.mesh.ip6addr='fe80::1/64'
EOF
