# src/packages/fff/fff-gateway/files/etc/uci-defaults/99-fff-gateway-sysupgrade
#!/bin/sh
# Copyright 2017 Tim Niemeyer
# License GPLv3

cat > /etc/sysupgrade.conf <<-__EOF__
/etc/config/gateway
__EOF__
