# src/packages/fff/fff-gateway/files/etc/uci-defaults/99-fff-gwupgrade
#!/bin/sh
# License GPLv3

cat >> /etc/sysupgrade.conf <<-__EOF__
/etc/config/gateway
/etc/hoodfile
__EOF__
