# src/packages/fff/fff-batman-adv-legacy-server/files/etc/uci-defaults/94-fff-batman-adv-server
#!/bin/sh

/bin/rm -f /usr/lib/micron.d/fff-batman-adv

uci batch <<EOF
  set batman-adv.bat0=mesh
  set batman-adv.bat0.gw_mode='server'
  set batman-adv.bat0.gw_sel_class='200000'
EOF

uci commit

exit 0
