# src/packages/fff/fff-babeld/files/etc/uci-defaults/02-network-rules
#!/bin/sh

uci batch <<EOF
  delete batman-adv.bat0
  add network rule
  set network.@rule[0].src='10.0.0.0/8'
  set network.@rule[0].lookup='10'
  set network.@rule[0].priority='20'
  set network.@rule[1]=rule
  set network.@rule[1].dest='10.0.0.0/8'
  set network.@rule[1].lookup='10'
  set network.@rule[1].priority='20'
  set network.@rule[2]=rule
  set network.@rule[2].in='mesh'
  set network.@rule[2].lookup='10'
  set network.@rule[2].priority='31'
EOF

uci commit

exit 0
